# systemd/misskey-cpuinfo.timer
[Unit]
Description=Run misskey-cpuinfo every 30 minutes

[Timer]
OnCalendar=*:00,30
Persistent=true

[Install]
WantedBy=timers.target
